# s6-overlay/s6-rc.d/conduwuit/run
#!/bin/sh
exec s6-setuidgid bitlbee /usr/local/bin/conduwuit \
  --config /var/lib/bitlbee/conduwuit/conduwuit.toml
