# s6-overlay/s6-rc.d/stunnel4/run
#!/bin/sh
exec s6-setuidgid bitlbee /usr/bin/stunnel4 /var/lib/bitlbee/stunnel.conf
