# s6-overlay/s6-rc.d/bitlbee/run
#!/bin/sh
exec s6-setuidgid bitlbee /usr/local/sbin/bitlbee -F -n -v -i 127.0.0.1
